# lighthouse-holesky.service

[Unit]
Description=Lighthouse Ethereum Client - Holesky Testnet
After=network.target geth-holesky.service
Wants=network.target geth-holesky.service

[Service]
User=root
Type=simple
ExecStart=/root/.cargo/bin/lighthouse bn \
  --network holesky \
  --datadir /root/holesky-node/data/lighthouse \
  --execution-endpoint http://localhost:8553 \
  --execution-jwt /root/holesky-node/data/jwtsecret \
  --checkpoint-sync-url https://checkpoint-sync.holesky.ethpandaops.io \
  --http \
  --http-address 0.0.0.0 \
  --http-port 5056 \
  --metrics \
  --metrics-address 0.0.0.0 \
  --metrics-port 5057 \
  --port 9002
Restart=on-failure
LimitNOFILE=65535

[Install]
WantedBy=default.target
